# docker/graphite-node/curator.cron
* * * * * root echo "Hello world" >> /var/log/cron.log 2>&1
# An empty line is required at the end of this file for a valid cron file.

#find /opt/graphite/storage/whisper/ -name "*wsp" -mtime +30 -delete
#find /opt/graphite/storage/whisper/ -type d -empty -delete

